# crontab
# m h dom mon dow command
0 * * * * sudo /home/riegel/scripts/backuphome.sh >> /home/riegel/log/crontab.log
